# FREEDESKTOP_MIME_TYPES_PATH/shared-mime-info-2.0/tests/mime-detection/systemd.automount
[Unit]
Description=Test Automount Point
Before=sysinit.target
ConditionPathExists=/some/mountpoint/

[Automount]
Where=/some/mountpoint
