# FREEDESKTOP_MIME_TYPES_PATH/shared-mime-info-2.0/tests/mime-detection/systemd.slice
[Unit]
Description=Cake Slice
Before=desserts.target
