# FREEDESKTOP_MIME_TYPES_PATH/shared-mime-info-2.0/tests/mime-detection/systemd.target
[Unit]
Description=Final Target
RefuseManualStart=yes
After=shutdown.target umount.target
