# FREEDESKTOP_MIME_TYPES_PATH/shared-mime-info-2.0/tests/mime-detection/systemd.socket
[Unit]
Description=Another Bus Messaging Socket

[Socket]
ListenStream=/run/another/system_bus_socket
