# FREEDESKTOP_MIME_TYPES_PATH/shared-mime-info-2.0/tests/mime-detection/systemd.service
[Unit]
Description=Pseudo-Emergency Shell
Conflicts=shutdown.target
Conflicts=rescue.service
Conflicts=syslog.socket
Before=shutdown.target

[Service]
Environment=HOME=/root
WorkingDirectory=-/root
Type=idle
