# FREEDESKTOP_MIME_TYPES_PATH/shared-mime-info-2.0/tests/mime-detection/systemd.path
[Unit]
Description=Some Kind of Path

[Path]
DirectoryNotEmpty=/run/systemd/some-path
MakeDirectory=yes
